# cron.d/gen_aprx_weather_packet
# generate (or beacon) weather packet for APRS every 10 mins.
# 
#
*/10   *     * * *     root  /usr/bin/perl /etc/scripts/aprx_weather_packet.pl
